# scripts/docquer.service
[Unit]
Description=Docquer - Homogeneous Torque

[Service]
Type=simple
ExecStart=/vagrant/docquer

[Install]
WantedBy=multi-user.target
