# scripts/docquer.service
[Unit]
Description=Docquer - Homogeneous Torque

[Service]
ExecStart=/vagrant/docquer
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
